# borg/.config/systemd/user/borg.timer
[Unit]
Description=Borg timer

[Timer]
OnCalendar=daily
Persistent=true
Unit=borg.service

[Install]
WantedBy=timers.target
