# systemd/.config/systemd/user/math-hw.timer
[Unit]
Description=Math Homework Daily Downloader

[Timer]
OnCalendar=daily
Persistent=true
Unit=math-hw.service

[Install]
WantedBy=timers.target
