# systemd/.config/systemd/user/moodle.timer
[Unit]
Description=Moodle Daily Downloader

[Timer]
OnCalendar=daily
Persistent=true
Unit=moodle.service

[Install]
WantedBy=timers.target
