# systemd/.config/systemd/user/moodle.service
[Unit]
Description=Automatic Moodle Donwloader

[Service]
Type=oneshot
ExecStart=%h/.local/bin/moodle-dl
WorkingDirectory=%h/docs/college/moodle/
