# tarsnap/.config/systemd/user/tarsnap.service
[Unit]
Description=Tarsnap Backup for docs

[Service]
Type=oneshot
ExecStart=%h/.config/tarsnap/tarsnap-backup.sh
