# tarsnap/.config/systemd/user/tarsnap.timer
[Unit]
Description=Tarsnap timer

[Timer]
OnCalendar=daily
Persistent=true
Unit=tarsnap.service

[Install]
WantedBy=timers.target
