# restic/.config/systemd/user/restic-backup.service
[Unit]
Description=Backup system

[Service]
ExecStart=systemd-inhibit %h/.config/restic/restic-backup.sh
